# deployments/proxmox/sluice.service
[Unit]
Description=Llama-CPP Sluice: Dynamic Asymmetric Inference
After=network.target

[Service]
Type=simple
User=skot
WorkingDirectory=/opt/swamc
Environment=SLUICE_MODEL_PATH=/models/qwen-3.6.gguf
Environment=SLUICE_TOTAL_POOL=98304
Environment=SLUICE_RESERVED_POOL=32768
ExecStart=/opt/swamc/venv/bin/python -m sluice.server
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
